# yt-sync.service
[Unit]
Description=YT Sync - YouTube Downloader Service
After=network.target

[Service]
Type=simple
User=__USER__
Group=__GROUP__
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/yt-sync-service
Restart=on-failure
RestartSec=10

# Environment
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yt-sync

[Install]
WantedBy=multi-user.target
